# scripts/saged.service
[Unit]
Description=SAGE Daemon - Process Manager
After=network.target

[Service]
WorkingDirectory=/home/ubuntu/projects/sage
ExecStart=/usr/local/bin/saged
Restart=always
Environment=HOME=/home/ubuntu
StandardOutput=file:/var/log/sage/saged.log
StandardError=file:/var/log/sage/saged-error.log

[Install]
WantedBy=multi-user.target
